# DocNodes/first-001-node/configs/doctor-first.service
[Unit]
Description=Doctor Admin - first-000-node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/doctor-admin
EnvironmentFile=/opt/doctor-admin/DocNodes/first-000-node/configs/first-000-node.env
ExecStart=/opt/doctor-admin/.venv/bin/python /opt/doctor-admin/main.py --env /opt/doctor-admin/DocNodes/first-000-node/configs/first-000-node.env
Restart=always
RestartSec=3
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
